# app-admin/dinit/autobuild/overrides/usr/lib/systemd/system/dinit.service
[Unit]
Description=Service manager

[Service]
ExecStart=/usr/bin/dinit --system

[Install]
WantedBy=default.target
