# app-admin/dinit/autobuild/overrides/usr/lib/systemd/user/dinit.service
[Unit]
Description=Service manager

[Service]
ExecStart=/usr/bin/dinit --user

[Install]
WantedBy=default.target
